# tmux/uas/guibox/supervisor.conf
[program:flight_summary]
command=/bin/bash %(ENV_KAMERA_DIR)s/tmux/uas/guibox/00_flight_summary.sh
user=root
autostart=true
