# tmux/uas/guibox/supervisor.conf
[program:flight_summary]
command=/bin/bash /home/user/kw/kamera/tmux/uas/guibox/00_flight_summary.sh
user=root
autostart=true
